# vscode/rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: VS Code Server
# Runs SSH daemon when enabled.
# ==============================================================================
if [[ -f /run/sshd.disabled ]]; then
    bashio::log.info "SSHD disabled; not starting sshd service"
    exec sleep infinity
fi

exec /usr/sbin/sshd -D -e -f /etc/ssh/sshd_config
